# etc/crontab
0 6 * * * /var/www/html/principal/Sinfin/etc/cron.sh daily
0 18 * * 7 /var/www/html/principal/Sinfin/etc/cron.sh weekly
